# AgroDrone-Edge-Node/services/agrodrone-mosaic.path
# etc/systemd/system/agrodrone-mosaic.path
[Unit]
Description=Trigger for processing the image data and sending it to the cloud

[Path]
PathChanged=/home/sr-design/system/flightplans/

[Install]
WantedBy=multi-user.target