# AgroDrone-Edge-Node/services/agrodrone-telemetry.path
# etc/systemd/system/agrodrone-telemetry.path
[Unit]
Description=Trigger for sending telemetry from Pi to Cloud

[Path]
PathChanged=/home/sr-design/system/telemetry.json

[Install]
WantedBy=multi-user.target