# AgroDrone-Edge-Node/services/agrodrone-flightplan-upload.path
# etc/systemd/system/agrodrone-flightplan-upload.path
[Unit]
Description=Trigger for sending the flight plan to the onboard Pi

[Path]
PathChanged=/home/sr-design/system/flightplans/waypoints.json

[Install]
WantedBy=multi-user.target